# deploy/systemd/led-clock.service
[Unit]
Description=MAX7219 LED Clock
After=network.target

[Service]
Type=simple
User=theza
WorkingDirectory=/opt/led-clock
EnvironmentFile=-/etc/default/led-clock
ExecStart=/opt/led-clock/.venv/bin/python /opt/led-clock/led_clock.py
Restart=always
RestartSec=2
StandardOutput=append:/var/log/led-clock.log
StandardError=append:/var/log/led-clock.log

[Install]
WantedBy=multi-user.target
